# sv/emacs/run
#!/bin/sh
exec chpst -u trajanus emacs -u trajanus --fg-daemon  2>&1
